# recipes-pem3/pem3-scripts/pem3-scripts/sdcard.rules
ACTION=="add",SUBSYSTEM=="block",KERNEL=="mmcblk1p1",RUN+="/usr/local/bin/software_update.sh"
